# deploy/r-file.service
[Unit]
Description=R-File 文件中转站后端服务
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/r-file/server
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
